# base/pkgsrc/root/var/lib/lrpkg/root.list
./
./etc/inittab
./etc/rc.d/rc.sysinit
./etc/coyote/
