# base/pkgsrc/acpid/var/lib/lrpkg/acpid.list
/sbin/acpid
/etc/acpi/*
/etc/rc.d/pkgs/rc.acpid
/var/lib/lrpkg/acpid.list
