# base/pkgsrc/advroute/var/lib/lrpkg/advroute.list
/lib/modules/advroute
/etc/modules.advroute
/etc/rc.d/pkgs/mod.advroute
/bin/vconfig
/var/lib/lrpkg/advroute.list
