# base/pkgsrc/modules/var/lib/lrpkg/modules.list
etc/modules
etc/rc.d/rc.modules
bin/vconfig
lib/modules
var/lib/lrpkg/modules.*
/etc/l7-protocols/*.pat
