# base/pkgsrc/modules/var/lib/lrpkg/modules.list
etc/modules
etc/rc.d/rc.modules
lib/modules
var/lib/lrpkg/modules.*
