# base/pkgsrc/l7filter/var/lib/lrpkg/l7filter.list
/etc/l7-protocols
/etc/rc.d/pkgs/mod.l7filter
/lib/modules/ipt_layer7.o
/var/lib/lrpkg/l7filter.list
