# base/pkgsrc/drivers/var/lib/lrpkg/drivers.list
/lib/modules/drivers
/etc/modules.drivers
/etc/rc.d/pkgs/mod.drivers
/var/lib/lrpkg/drivers.list
