# base/pkgsrc/bridge/var/lib/lrpkg/bridge.list
/lib/libebt*.so
/lib/modules/bridge
/etc/modules.bridge
/bin/brctl
/bin/ebtables
/etc/ethertypes
/etc/rc.d/pkgs/mod.zz_bridge
/var/lib/lrpkg/bridge.list
