# base/pkgsrc/webadmin/var/lib/lrpkg/webadmin.list
/etc/rc.d/pkgs/rc.webadmin
/etc/lighttpd.tpl
/usr/lib/libpcre.so.0
/usr/lib/libpcre.so
/usr/lib/libpcreposix.so.0
/usr/lib/libpcreposix.so
/usr/lib/libpcrecpp.so.0
/usr/lib/libpcrecpp.so
/usr/lib/libpcre.so.0.0.1
/usr/lib/libpcreposix.so.0.0.0
/usr/lib/libpcrecpp.so.0.0.0
/usr/lib/mod_access.so
/usr/lib/mod_accesslog.so
/usr/lib/mod_alias.so
/usr/lib/mod_auth.so
/usr/lib/mod_cgi.so
/usr/lib/mod_compress.so
/usr/lib/mod_dirlisting.so
/usr/lib/mod_fastcgi.so
/usr/lib/mod_flv_streaming.so
/usr/lib/mod_indexfile.so
/usr/lib/mod_scgi.so
/usr/lib/mod_staticfile.so
/usr/bin/proccgi
/usr/bin/cryptpw
/usr/bin/pcretest
/usr/bin/pcregrep
/usr/sbin/thttpd
/usr/sbin/lighttpd
/var/http
/var/lib/lrpkg/webadmin.*