# base/pkgsrc/webadmin/var/lib/lrpkg/webadmin.list
/var/lib/lrpkg/webadmin.*
/etc/rc.d/pkgs/rc.webadmin
/etc/lighttpd.tpl
/var/http/header.inc
/var/http/style.css
/var/http/web-functions
/var/http/htdocs/brasao.png
/var/http/htdocs/logo.png
/var/http/htdocs/index.htm
/var/http/htdocs/cgi-bin/qosclasses.cgi
/var/http/htdocs/cgi-bin/inetconf.cgi
/var/http/htdocs/cgi-bin/hangup-ppp.cgi
/var/http/htdocs/cgi-bin/qosfilter.cgi
/var/http/htdocs/cgi-bin/qos.cgi
/var/http/htdocs/cgi-bin/editconf.cgi
/var/http/htdocs/cgi-bin/profiles.cgi
/var/http/htdocs/cgi-bin/lanconf.cgi
/var/http/htdocs/cgi-bin/reboot.cgi
/var/http/htdocs/cgi-bin/index.cgi
/var/http/htdocs/cgi-bin/chpasswd.cgi
/var/http/htdocs/cgi-bin/dial-ppp.cgi
/var/http/htdocs/cgi-bin/diags.cgi
/var/http/htdocs/cgi-bin/qosstatus.cgi
/var/http/htdocs/cgi-bin/backup.cgi
/var/http/htdocs/cgi-bin/subnet.cgi
/var/http/htdocs/cgi-bin/portfw.cgi
/var/http/htdocs/cgi-bin/localhosts.cgi
/var/http/htdocs/cgi-bin/loadbalance.cgi
/var/http/htdocs/cgi-bin/leases.cgi
/var/http/htdocs/cgi-bin/firewall.cgi
/var/http/htdocs/cgi-bin/dhcpconf.cgi
/var/http/htdocs/cgi-bin/cron.cgi
/var/http/htdocs/cgi-bin/adminconf.cgi
/var/http/htdocs/cgi-bin/sfirewall.cgi
/usr/bin/cryptpw
/usr/bin/pcregrep
/usr/bin/pcretest
/usr/sbin/lighttpd
/usr/sbin/thttpd
/usr/bin/proccgi
/usr/lib/libpcre.so
/usr/lib/libpcre.so.0
/usr/lib/libpcre.so.0.0.1
/usr/lib/libpcrecpp.so
/usr/lib/libpcrecpp.so.0
/usr/lib/libpcrecpp.so.0.0.0
/usr/lib/libpcreposix.so
/usr/lib/libpcreposix.so.0
/usr/lib/libpcreposix.so.0.0.0
/usr/lib/mod_access.so
/usr/lib/mod_accesslog.so
/usr/lib/mod_alias.so
/usr/lib/mod_auth.so
/usr/lib/mod_cgi.so
/usr/lib/mod_compress.so
/usr/lib/mod_dirlisting.so
/usr/lib/mod_fastcgi.so
/usr/lib/mod_flv_streaming.so
/usr/lib/mod_indexfile.so
/usr/lib/mod_scgi.so
/usr/lib/mod_staticfile.so
/usr/lib/mod_status.so
/var/http/htdocs/imagens_doacao/*
