# base/pkgsrc/webadmin/var/lib/lrpkg/webadmin.list
/etc/rc.d/pkgs/rc.webadmin
/usr/sbin/thttpd
/usr/bin/proccgi
/usr/bin/cryptpw
/var/lib/lrpkg/webadmin.*
/var/http
