# base/pkgsrc/modules_smp/var/lib/lrpkg/modules.list
/etc/modules
/etc/rc.d/rc.modules
/etc/rc.d/pkgs/mod.zz_bridge
/etc/modules.bridge
/etc/l7-protocols/radmin.pat
/etc/l7-protocols/telnet.pat
/etc/l7-protocols/ssh.pat
/etc/l7-protocols/cvs.pat
/etc/l7-protocols/http.pat
/etc/l7-protocols/halflife2-deathmatch.pat
/etc/l7-protocols/chikka.pat
/etc/l7-protocols/gnucleuslan.pat
/etc/l7-protocols/code_red.pat
/etc/l7-protocols/snmp-mon.pat
/etc/l7-protocols/goboogy.pat
/etc/l7-protocols/pressplay.pat
/etc/l7-protocols/skypetoskype.pat
/etc/l7-protocols/liveforspeed.pat
/etc/l7-protocols/uucp.pat
/etc/l7-protocols/applejuice.pat
/etc/l7-protocols/whois.pat
/etc/l7-protocols/audiogalaxy.pat
/etc/l7-protocols/rtf.pat
/etc/l7-protocols/100bao.pat
/etc/l7-protocols/qq.pat
/etc/l7-protocols/ares.pat
/etc/l7-protocols/ipp.pat
/etc/l7-protocols/soribada.pat
/etc/l7-protocols/postscript.pat
/etc/l7-protocols/x11.pat
/etc/l7-protocols/rtsp.pat
/etc/l7-protocols/socks.pat
/etc/l7-protocols/xboxlive.pat
/etc/l7-protocols/quake-halflife.pat
/etc/l7-protocols/ident.pat
/etc/l7-protocols/httpvideo.pat
/etc/l7-protocols/vnc.pat
/etc/l7-protocols/tor.pat
/etc/l7-protocols/tesla.pat
/etc/l7-protocols/bgp.pat
/etc/l7-protocols/ssl.pat
/etc/l7-protocols/subspace.pat
/etc/l7-protocols/ciscovpn.pat
/etc/l7-protocols/tsp.pat
/etc/l7-protocols/http-dap.pat
/etc/l7-protocols/battlefield2.pat
/etc/l7-protocols/ftp.pat
/etc/l7-protocols/gtalk.pat
/etc/l7-protocols/skypeout.pat
/etc/l7-protocols/ogg.pat
/etc/l7-protocols/napster.pat
/etc/l7-protocols/http-freshdownload.pat
/etc/l7-protocols/lpd.pat
/etc/l7-protocols/guildwars.pat
/etc/l7-protocols/counterstrike-source.pat
/etc/l7-protocols/battlefield2142.pat
/etc/l7-protocols/pdf.pat
/etc/l7-protocols/validcertssl.pat
/etc/l7-protocols/stun.pat
/etc/l7-protocols/snmp.pat
/etc/l7-protocols/nimda.pat
/etc/l7-protocols/soulseek.pat
/etc/l7-protocols/shoutcast.pat
/etc/l7-protocols/replaytv-ivs.pat
/etc/l7-protocols/gnutella.pat
/etc/l7-protocols/msn-filetransfer.pat
/etc/l7-protocols/hddtemp.pat
/etc/l7-protocols/doom3.pat
/etc/l7-protocols/hotline.pat
/etc/l7-protocols/exe.pat
/etc/l7-protocols/zmaap.pat
/etc/l7-protocols/teamspeak.pat
/etc/l7-protocols/imap.pat
/etc/l7-protocols/nbns.pat
/etc/l7-protocols/smb.pat
/etc/l7-protocols/httpcachehit.pat
/etc/l7-protocols/bittorrent.pat
/etc/l7-protocols/battlefield1942.pat
/etc/l7-protocols/fasttrack.pat
/etc/l7-protocols/tftp.pat
/etc/l7-protocols/citrix.pat
/etc/l7-protocols/http-rtsp.pat
/etc/l7-protocols/jpeg.pat
/etc/l7-protocols/mp3.pat
/etc/l7-protocols/zip.pat
/etc/l7-protocols/quake1.pat
/etc/l7-protocols/httpcachemiss.pat
/etc/l7-protocols/pop3.pat
/etc/l7-protocols/gopher.pat
/etc/l7-protocols/freenet.pat
/etc/l7-protocols/cimd.pat
/etc/l7-protocols/rtp.pat
/etc/l7-protocols/worldofwarcraft.pat
/etc/l7-protocols/tonghuashun.pat
/etc/l7-protocols/msnmessenger.pat
/etc/l7-protocols/thecircle.pat
/etc/l7-protocols/yahoo.pat
/etc/l7-protocols/finger.pat
/etc/l7-protocols/kugoo.pat
/etc/l7-protocols/quicktime.pat
/etc/l7-protocols/smtp.pat
/etc/l7-protocols/flash.pat
/etc/l7-protocols/jabber.pat
/etc/l7-protocols/ncp.pat
/etc/l7-protocols/http-itunes.pat
/etc/l7-protocols/dns.pat
/etc/l7-protocols/h323.pat
/etc/l7-protocols/rpm.pat
/etc/l7-protocols/directconnect.pat
/etc/l7-protocols/imesh.pat
/etc/l7-protocols/poco.pat
/etc/l7-protocols/openft.pat
/etc/l7-protocols/aim.pat
/etc/l7-protocols/edonkey.pat
/etc/l7-protocols/dayofdefeat-source.pat
/etc/l7-protocols/gif.pat
/etc/l7-protocols/unknown.pat
/etc/l7-protocols/perl.pat
/etc/l7-protocols/unset.pat
/etc/l7-protocols/runesofmagic.pat
/etc/l7-protocols/biff.pat
/etc/l7-protocols/xunlei.pat
/etc/l7-protocols/mute.pat
/etc/l7-protocols/pplive.pat
/etc/l7-protocols/dhcp.pat
/etc/l7-protocols/gkrellm.pat
/etc/l7-protocols/httpaudio.pat
/etc/l7-protocols/irc.pat
/etc/l7-protocols/sip.pat
/etc/l7-protocols/dazhihui.pat
/etc/l7-protocols/rdp.pat
/etc/l7-protocols/png.pat
/etc/l7-protocols/rlogin.pat
/etc/l7-protocols/nntp.pat
/etc/l7-protocols/html.pat
/etc/l7-protocols/tar.pat
/etc/l7-protocols/netbios.pat
/etc/l7-protocols/ntp.pat
/etc/l7-protocols/snmp-trap.pat
/etc/l7-protocols/mohaa.pat
/etc/l7-protocols/subversion.pat
/etc/l7-protocols/ventrilo.pat
/etc/l7-protocols/aimwebcontent.pat
/etc/l7-protocols/teamfortress2.pat
/etc/l7-protocols/armagetron.pat
/etc/l7-protocols/ssdp.pat
/etc/l7-protocols/live365.pat
/etc/l7-protocols/rar.pat
/etc/l7-protocols/pcanywhere.pat
/lib/modules/drivers/parport/parport.o
/lib/modules/drivers/parport/parport_pc.o
/lib/modules/drivers/usb/usbnet.o
/lib/modules/drivers/usb/pegasus.o
/lib/modules/drivers/usb/catc.o
/lib/modules/drivers/usb/CDCEther.o
/lib/modules/drivers/usb/rtl8150.o
/lib/modules/drivers/usb/acm.o
/lib/modules/drivers/usb/serial/usbserial.o
/lib/modules/drivers/usb/serial/pl2303.o
/lib/modules/drivers/usb/kaweth.o
/lib/modules/drivers/net/r8169.o
/lib/modules/drivers/net/fealnx.o
/lib/modules/drivers/net/dmfe.o
/lib/modules/drivers/net/tg3.o
/lib/modules/drivers/net/amd8111e.o
/lib/modules/drivers/net/hamachi.o
/lib/modules/drivers/net/yellowfin.o
/lib/modules/drivers/net/acenic.o
/lib/modules/drivers/net/mii.o
/lib/modules/drivers/net/e1000/e1000.o
/lib/modules/drivers/net/tulip/tulip.o
/lib/modules/drivers/net/sk98lin/sk98lin.o
/lib/modules/drivers/net/sundance.o
/lib/modules/drivers/net/velocityget/velocityget.o
/lib/modules/drivers/net/8390.o
/lib/modules/drivers/net/pcnet32.o
/lib/modules/drivers/net/wireless/hermes.o
/lib/modules/drivers/net/wireless/orinoco.o
/lib/modules/drivers/net/ns83820.o
/lib/modules/drivers/net/via-rhine.o
/lib/modules/drivers/net/tlan.o
/lib/modules/drivers/net/eepro100.o
/lib/modules/drivers/net/8139cp.o
/lib/modules/drivers/net/forcedeth.o
/lib/modules/drivers/net/skge.o
/lib/modules/drivers/net/de4x5.o
/lib/modules/drivers/net/ne2k-pci.o
/lib/modules/drivers/net/8139too.o
/lib/modules/drivers/net/winbond-840.o
/lib/modules/drivers/net/e100/e100.o
/lib/modules/drivers/net/3c59x.o
/lib/modules/drivers/net/ppp_synctty.o
/lib/modules/drivers/net/sis900.o
/lib/modules/drivers/net/typhoon.o
/lib/modules/drivers/net/bsd_comp.o
/lib/modules/drivers/net/ppp_async.o
/lib/modules/drivers/net/sky2.o
/lib/modules/drivers/net/dgrs.o
/lib/modules/drivers/net/slip.o
/lib/modules/drivers/net/bonding/bonding.o
/lib/modules/drivers/net/tun.o
/lib/modules/drivers/net/dl2k.o
/lib/modules/drivers/net/ppp_mppe_mppc.o
/lib/modules/drivers/net/b44.o
/lib/modules/drivers/net/natsemi.o
/lib/modules/drivers/net/starfire.o
/lib/modules/drivers/net/hp100.o
/lib/modules/drivers/net/ppp_deflate.o
/lib/modules/drivers/net/epic100.o
/lib/modules/drivers/char/parport_serial.o
/lib/modules/lib/zlib_inflate/zlib_inflate.o
/lib/modules/lib/zlib_deflate/zlib_deflate.o
/lib/modules/net/ipv4/ipip.o
/lib/modules/net/ipv4/netfilter/ipt_recent.o
/lib/modules/net/ipv4/netfilter/ipt_ecn.o
/lib/modules/net/ipv4/netfilter/ipt_ULOG.o
/lib/modules/net/ipv4/netfilter/ipt_owner.o
/lib/modules/net/ipv4/netfilter/ip_conntrack_amanda.o
/lib/modules/net/ipv4/netfilter/ip_queue.o
/lib/modules/net/ipv4/netfilter/ipt_ECN.o
/lib/modules/net/ipv4/netfilter/ip_nat_irc.o
/lib/modules/net/ipv4/netfilter/ipt_layer7.o
/lib/modules/net/ipv4/netfilter/ipt_ah.o
/lib/modules/net/ipv4/netfilter/ip_conntrack_tftp.o
/lib/modules/net/ipv4/netfilter/ip_nat_amanda.o
/lib/modules/net/ipv4/netfilter/ipt_MIRROR.o
/lib/modules/net/ipv4/netfilter/ipt_esp.o
/lib/modules/net/ipv4/netfilter/ipt_dscp.o
/lib/modules/net/ipv4/netfilter/ip_conntrack_irc.o
/lib/modules/net/ipv4/netfilter/ipt_DSCP.o
/lib/modules/net/ipv4/netfilter/ip_nat_tftp.o
/lib/modules/net/ipv4/ip_gre.o
/lib/modules/net/8021q/8021q.o
/lib/modules/net/sctp/sctp.o
/lib/modules/net/bridge/netfilter/ebt_limit.o
/lib/modules/net/bridge/netfilter/ebt_stp.o
/lib/modules/net/bridge/netfilter/ebt_vlan.o
/lib/modules/net/bridge/netfilter/ebt_mark_m.o
/lib/modules/net/bridge/netfilter/ebt_802_3.o
/lib/modules/net/bridge/netfilter/ebtables.o
/lib/modules/net/bridge/netfilter/ebt_log.o
/lib/modules/net/bridge/netfilter/ebt_redirect.o
/lib/modules/net/bridge/netfilter/ebt_ulog.o
/lib/modules/net/bridge/netfilter/ebtable_broute.o
/lib/modules/net/bridge/netfilter/ebt_arpreply.o
/lib/modules/net/bridge/netfilter/ebt_pkttype.o
/lib/modules/net/bridge/netfilter/ebt_dnat.o
/lib/modules/net/bridge/netfilter/ebt_snat.o
/lib/modules/net/bridge/netfilter/ebt_ip.o
/lib/modules/net/bridge/netfilter/ebtable_filter.o
/lib/modules/net/bridge/netfilter/ebt_mark.o
/lib/modules/net/bridge/netfilter/ebt_arp.o
/lib/modules/net/bridge/netfilter/ebt_among.o
/lib/modules/net/bridge/netfilter/ebtable_nat.o
/lib/modules/net/bridge/bridge.o
/lib/modules/net/sched/sch_hfsc.o
/lib/modules/net/sched/sch_gred.o
/lib/modules/net/sched/sch_csz.o
/lib/modules/net/sched/sch_cbq.o
/lib/modules/net/sched/sch_red.o
/lib/modules/net/sched/sch_teql.o
/lib/modules/net/sched/sch_tbf.o
/lib/modules/net/sched/cls_rsvp.o
/lib/modules/net/sched/sch_dsmark.o
/lib/modules/crypto/arc4.o
/lib/modules/crypto/des.o
/lib/modules/crypto/md4.o
/lib/modules/crypto/sha512.o
/lib/modules/crypto/sha256.o
/lib/modules/crypto/deflate.o
/lib/modules/crypto/aes.o
/lib/modules/crypto/sha1.o
/bin/vconfig
/var/lib/lrpkg/modules.list
