# base/pkgsrc/hdparm/var/lib/lrpkg/hdparm.list
/sbin/hdparm
/etc/rc.d/pkgs/mod.hdparm
/var/lib/lrpkg/hdparm.list